# etc/supervisord.d/carbon-relay-ng_conf.ini
[program:carbon-relay-ng-conf]
command=consul-template -consul localhost:8500 -wait 5s -template "/etc/consul-templates/carbon-relay-ng.ini.ctmpl:/etc/carbon-relay-ng/carbon-relay-ng.ini:supervisorctl restart carbon-relay-ng"
redirect_stderr=true
autostart=true
startsecs=10
stdout_logfile=/var/log/supervisor/%(program_name)s.log
