# etc/supervisord.d/carbon-relay-ng.ini
[program:carbon-relay-ng]
command = carbon-relay-ng  etc/carbon-relay-ng/carbon-relay-ng.ini
redirect_stderr=true
stdout_logfile=syslog
autostart = false
